# deployments/ubuntu/etc/cron.d/amsco_ssl_cert_renewal
MAILTO="tino@erpweb.co.za"
0 0 1 * * /var/letsencrypt/./letsencrypt-auto certonly -a webroot --webroot-path=/usr/share/nginx/html -d sbsc.serthe.com