# deployments/ubuntu/etc/cron.d/ssl_cert_renewal
# m h  dom mon dow  command
MAILTO="brad@breadcrumb.tech"
0 5 1 * * /var/certbot/certbot-auto renew --quiet --no-self-upgrade  >> /var/log/letsencrypt-renewal.log
